# board/rootfs_overlay/etc/systemd/system/mount-application-partition.service
[Unit]
Description=Test hello-world based service

[Service]
ExecStart=/usr/sbin/mounts/rootfs-custom-mount.sh

[Install]
WantedBy=multi-user.target